# TheLoom/the-loom/the-loom.service
# TheLoom - Hidden State Extraction Server
# Part of the Weaver Ecosystem
#
# Installation:
#   1. Copy this file: sudo cp the-loom.service /etc/systemd/system/
#   2. Edit the service file to customize for your system:
#      sudo nano /etc/systemd/system/the-loom.service
#      - Update User/Group (lines 25-26) to your username
#      - Update WorkingDirectory (line 29) to your TheLoom install path
#      - Update ExecStart (line 30) to your poetry path
#   3. Reload systemd: sudo systemctl daemon-reload
#   4. Enable and start: sudo systemctl enable --now the-loom
#
# View logs:
#   journalctl -u the-loom -f
#
# Check status:
#   systemctl status the-loom

[Unit]
Description=TheLoom Hidden State Extraction Server
Documentation=https://github.com/r3d91ll/WeaverTools
After=network.target

[Service]
Type=simple
# TODO: Change to your username
User=YOUR_USERNAME
Group=YOUR_USERNAME

# TODO: Change to your TheLoom installation directory
WorkingDirectory=/path/to/WeaverTools/TheLoom/the-loom
# TODO: Change to your poetry path (run 'which poetry' to find it)
ExecStart=/path/to/poetry run loom --port 8080

# Restart behavior
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
Environment=LOOM_LOGGING__LEVEL=INFO

# Resource limits (adjust as needed for your GPU memory)
# MemoryMax=32G
# CPUQuota=400%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
